# root/etc/services.d/monitor/run
#!/usr/bin/with-contenv bash

# shellcheck shell=bash

while [ ! -f /var/run/vpn-up ]; do
    sleep 1
done

while [ -f /var/run/vpn-up ]; do
    sleep 1
done
