# root/etc/services.d/openvpn/finish
#!/usr/bin/execlineb -S0

# shellcheck shell=bash

# This kills the container
s6-svscanctl -t /var/run/s6/services
